# bin/build.cron
#!/bin/sh
exec runexclusive -f ~/tmp/build.lock --no-wait -t 36000 build
